# supervisor/baraka-ai-dual.conf
[program:vllm]
command=/workspace/start_vllm_dual.sh
autostart=true
autorestart=true
startretries=5
startsecs=30
stdout_logfile=/workspace/vllm.log
stderr_logfile=/workspace/vllm.err.log

[program:ollama]
command=ollama serve
autostart=true
autorestart=true
startretries=5
startsecs=10
stdout_logfile=/workspace/ollama.log
stderr_logfile=/workspace/ollama.err.log
environment=OLLAMA_HOST="0.0.0.0:11434",OLLAMA_ORIGINS="*",OLLAMA_KEEP_ALIVE="24h",CUDA_VISIBLE_DEVICES="",GGML_VK_VISIBLE_DEVICES=""
